# setup/foundry/foundry_healthcheck.service
[Unit]
Description=Foundry VTT Health Check

[Service]
Type=oneshot
ExecStart=/bin/bash /foundrycron/foundry_healthcheck.sh
